# ipwatcher.service

[Unit]
Description=IpWatcher service
After=network.target

[Service]
ExecStart=/home/ubuntu/pgsm/ipwatcher.sh

[Install]
WantedBy=multi-user.target
